# rootfs/etc/s6-overlay/s6-rc.d/aria2-log/run
#!/bin/sh
exec logutil-service /var/log/aria2
