# rootfs/etc/s6-overlay/s6-rc.d/smartdns/run
#!/bin/sh
. /bin/log_functions.sh

# SmartDNS服务启动脚本
# 使用s6-setuidgid以非root用户启动，通过s6-notifyoncheck检测就绪状态
SMARTDNS_PORT="${SMARTDNS_PORT:-5353}"

log_info "smartdns" "正在启动 SmartDNS 服务 (端口:${SMARTDNS_PORT})..."

# 就绪检测：通过TCP连接检测SmartDNS端口是否已监听
exec 2>&1
exec s6-notifyoncheck -d -w 5000 -n 10 \
    -c "nc -w 1 127.0.0.1 ${SMARTDNS_PORT} < /dev/null > /dev/null 2>&1" \
    -- s6-setuidgid ddsuser smartdns -c /config/smartdns/smartdns.conf
