# rootfs/etc/s6-overlay/s6-rc.d/transmission-log/run
#!/bin/sh
exec logutil-service /var/log/transmission
