# rootfs/etc/s6-overlay/s6-rc.d/qbittorrent/run
#!/bin/sh
. /bin/log_functions.sh

# qBittorrent服务启动脚本
# 使用s6-setuidgid以非root用户启动
# 使用s6-notifyoncheck检测WebUI就绪状态
QBITTORRENT_WEBUI_PORT="${QBITTORRENT_WEBUI_PORT:-8080}"

log_info "qbittorrent" "正在启动 qBittorrent 服务 (WebUI端口:${QBITTORRENT_WEBUI_PORT})..."

# 就绪检测：通过HTTP请求检测WebUI是否可访问
exec 2>&1
exec s6-notifyoncheck -d -w 5000 -n 10 \
    -c "curl -sf http://localhost:${QBITTORRENT_WEBUI_PORT}/ > /dev/null 2>&1" \
    -- s6-setuidgid ddsuser qbittorrent-nox --profile=/config/qbittorrent --webui-port=${QBITTORRENT_WEBUI_PORT}
