# rootfs/etc/s6-overlay/s6-rc.d/adguardhome/finish
#!/bin/sh
/bin/s6-finish-handler.sh adguardhome "$1"
