# rootfs/etc/s6-overlay/s6-rc.d/transmission/finish
#!/bin/sh
/bin/s6-finish-handler.sh transmission "$1"
