# rootfs/etc/s6-overlay/s6-rc.d/transmission/run
#!/bin/sh
. /bin/log_functions.sh

# Transmission服务启动脚本
# 使用s6-setuidgid以非root用户启动
# 使用s6-notifyoncheck检测RPC接口就绪状态
TRANSMISSION_RPC_PORT="${TRANSMISSION_RPC_PORT:-9091}"

log_info "transmission" "正在启动 Transmission 服务 (RPC端口:${TRANSMISSION_RPC_PORT})..."

# 就绪检测：通过HTTP请求检测RPC接口是否可访问
exec 2>&1
exec s6-notifyoncheck -d -w 5000 -n 10 \
    -c "curl -sf http://localhost:${TRANSMISSION_RPC_PORT}/transmission/rpc > /dev/null 2>&1 || curl -sf http://localhost:${TRANSMISSION_RPC_PORT}/ > /dev/null 2>&1" \
    -- s6-setuidgid ddsuser transmission-daemon --foreground --config-dir /config/transmission --log-error
